# etc/service/apache/run
#!/bin/sh
rm -f /var/lock/httpd.pid
exec /usr/sbin/apache2ctl -D FOREGROUND
